# infra/systemd-user/oc-dev.service
[Unit]
Description=OpenClaw dev agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/openclawuser/.openclaw/workspace
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/openclawuser/.openclaw/workspace/scripts/oc-role-runner.sh dev
Restart=on-failure
RestartSec=5
CPUQuota=160%
MemoryMax=2200M
TasksMax=512
TimeoutStopSec=20

[Install]
WantedBy=default.target
