# infra/systemd-user/oc-youtube.service
[Unit]
Description=OpenClaw youtube agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/openclawuser/.openclaw/workspace
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/openclawuser/.openclaw/workspace/scripts/oc-role-runner.sh youtube
Restart=on-failure
RestartSec=5
CPUQuota=100%
MemoryMax=1200M
TasksMax=512
TimeoutStopSec=20

[Install]
WantedBy=default.target
