# infra/systemd-user/pay2pass-bridge.timer
[Unit]
Description=Pay2Pass bridge daily expiry timer at 00:00 UTC

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true
Unit=pay2pass-bridge-expiry.service

[Install]
WantedBy=timers.target
